# deployment/api-production.service
[Unit]
Description=Madrasha Backend API - Production
After=network.target

[Service]
WorkingDirectory=/var/www/api-production
ExecStart=/root/.bun/bin/bun run start
Restart=always
RestartSec=5

# Environment
Environment=NODE_ENV=production
Environment=PORT=5000

# Load environment variables from file
EnvironmentFile=/var/www/api-production/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/www/api-production
ReadWritePaths=/var/www/api-production/logs

[Install]
WantedBy=multi-user.target